# conf/dev/docker.elasticsearch.redis.ini
[program:docker-elasticsearch-redis]
directory = /home/ubuntu/poc/vue-storefront-api
user = ubuntu
command = docker-compose up
stdout_logfile = /home/ubuntu/poc/log/docker-elasticsearch-redis-out.log
stderr_logfile = /home/ubuntu/poc/log/docker-elasticsearch-redis-err.log
autostart = true
autorestart = true
killasgroup = true
stopasgroup = true
startretries = 5
startsecs = 0
stopwaitsecs = 600
numprocs = 1