# conf/dev/pocdjvue.gunicorn.ini
[program:pocdjvue-gunicorn]
directory = /home/ubuntu/poc/pocdjvue
user = ubuntu
command = /home/ubuntu/poc/pocdjvue/conf/dev/pocdjvue.gunicorn.sh
stdout_logfile = /home/ubuntu/poc/log/pocdjvue.gunicorn-out.log
stderr_logfile = /home/ubuntu/poc/log/pocdjvue-gunicorn-err.log
pidfile = /home/ubuntu/poc/pocdjvue.pid
autostart = true
autorestart = true
killasgroup = true
startretries = 5
startsecs = 0
stopwaitsecs = 600
numprocs = 1