# conf/dev/vue-storefront.ini
[program:vue-storefront]
directory = /home/ubuntu/poc/vue-storefront
user = ubuntu
command = yarn dev
stdout_logfile = /home/ubuntu/poc/log/vue-storefront-out.log
stderr_logfile = /home/ubuntu/poc/log/vue-storefront-err.log
autostart = true
autorestart = true
killasgroup = true
stopasgroup = true
startretries = 5
startsecs = 0
stopwaitsecs = 600
numprocs = 1